# hybris-quirks/lmk-disable.initd
#!/sbin/openrc-run

depend() {
    need lxc-android
}

start() {
    /usr/bin/lmk-disable.sh
}

stop() {
    :
}
