# lxc-android/etc/init.d/hybrislinks
#!/sbin/openrc-run

depend() {
    after android-mount lxc-android
    before udevhack bluebinder phosh ofono ofono2mm hadess-sensorfw-proxy audiosystem-passthrough
}

start() {
    ebegin "Starting hybrislinks service"
    link_libc
    eend $?
}

stop() {
    ebegin "Stopping hybrislinks service"
    eend 0
}

link_libc() {
    LXC_ROOTFS="/var/lib/lxc/android/rootfs"
    ANDROID_SDK_VERSION=$(grep ro.build.version.sdk= ${LXC_ROOTFS}/system/build.prop | cut -d "=" -f2)

    if [ "$(getconf LONG_BIT)" = "32" ]; then
        LIBDIR="lib"
    else
        LIBDIR="lib64"
    fi

    if [ -f "/android/vendor/$LIBDIR/egl/libGLESv2_adreno.so" ]; then
         ln -sf /android/vendor/$LIBDIR/egl/libGLESv2_adreno.so /libGLESv2_adreno.so
    fi

    if [ -f "/android/vendor/$LIBDIR/egl/eglSubDriverAndroid.so" ]; then
         ln -sf /android/vendor/$LIBDIR/egl/eglSubDriverAndroid.so /eglSubDriverAndroid.so
    fi

    if [ ! -f "/android/system/$LIBDIR/libaudioroute.so" ]; then
         ln -sf /android/system/$LIBDIR/vndk-$ANDROID_SDK_VERSION/libaudioroute.so /libaudioroute.so
    elif [ -f "/android/system/$LIBDIR/libaudioroute.so" ]; then
         ln -sf /android/system/$LIBDIR/libaudioroute.so /libaudioroute.so
    fi

    if [ $ANDROID_SDK_VERSION -ge 29 ]; then

        dir_list="/android/system/apex/com.android.runtime/${LIBDIR}/bionic
/android/system/apex/com.android.runtime.release/${LIBDIR}/bionic
/android/system/apex/com.android.runtime.debug/${LIBDIR}/bionic"

        lib_list="libc.so libm.so libdl.so"

        for dir in $dir_list; do
            if [ -d "$dir" ]; then
                for lib in $lib_list; do
                    if [ -f "${dir}/${lib}" ]; then
                        ln -sf "${dir}/${lib}" "/" || true
                        echo "Symlink created: ${dir}/${lib} -> /${lib}"
                    else
                        echo "Warning: ${dir}/${lib} does not exist."
                    fi
                done
                break
            fi
        done
    else
        true
    fi
}
